# arch/_to_sort/home/.config/systemd/user/battery-alert.timer
[Unit]
Description=Runs battery alert script periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
