# arch/home/.config/systemd/user/battery-alert.service
[Unit]
Description=Battery Level Notification

[Service]
ExecStart=/home/phantomwise/.local/bin/scripts/battery-alert.sh
