# arch/home/.config/systemd/user/battery-alert.service
[Unit]
Description=Battery Level Notification

[Service]
ExecStart=/home/phantomwise/.config/dunst/battery-alert.sh
